# environments/dev/sources/ec2/user_data/comfyui/comfyui.service
[Unit]
Description=ComfyUI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ComfyUI
ExecStart=/bin/bash -lc 'source .venv/bin/activate && python3 main.py --listen'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
